# templates/ld.so.preload.j2
{{ snoopy_ld_preload if snoopy_ld_preload is defined else '/usr/$LIB/libsnoopy.so' }}
