# templates/ld.so.preload.j2
{{ snoopy_ld_preload | default('/usr/$LIB/libsnoopy.so') }}
